# supervisord.conf
[supervisord]
nodaemon = true

[program:icecast]
command = /usr/local/bin/icecast -c /icecast/icecast.xml
user = root
autostart = true
